# etc/systemd/system/x11vnc.service
[Unit]
After=xdm.service

[Service]
ExecStart=/usr/local/bin/startx11vnc.sh

[Install]
WantedBy=default.target
